# etc/s6-overlay/s6-rc.d/vsftpd/run
#!/command/execlineb -P
/usr/local/bin/init-vsftpd