# etc/s6-overlay/s6-rc.d/log-tail/run
#!/command/execlineb -P
tail -f /var/log/vsftpd.log